# git/etc/profile.d/git-prompt.sh.REMOVED.git-id
9f8a4341a3bf6bd735abe13e348fad4e4d7915a2